# ansible/roles/nginx_letsencrypt_rproxy/files/certbot-renew.service
[Unit]
Description=Renew certbot certificates
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/certbot/certbot-auto renew
